# script/auto-dev/com.helix.auto-dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.helix.auto-dev</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/onetwo/.npm-global/bin/bun</string>
        <string>run</string>
        <string>/Users/onetwo/Documents/trae_projects/Helix/script/auto-dev/scheduler.ts</string>
        <string>--once</string>
        <string>--chat-id</string>
        <string>oc_adb8a26e34983ef16da7589560e24f5a</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/onetwo/Documents/trae_projects/Helix</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>14</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/onetwo/.local/share/mimocode/log/auto-dev-launchd.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/onetwo/.local/share/mimocode/log/auto-dev-launchd-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/onetwo/.npm-global/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/onetwo</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
